# bms.monitor.pi/bms-monitor-pi.service
[Unit]
Description=BMS BLE central daemon (JK BD6A24S10P family)
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
# Every start (including systemd auto-restart after a crash) does a
# fresh bluetoothd bounce. BlueZ regularly wedges into an "InProgress"
# state after our earlier session left a half-open connection, and no
# amount of Python-side retry recovers from that; a service restart
# clears it. Passwordless sudo for `systemctl restart bluetooth` is
# already granted to the pi user on this Pi (default in bluez.rules).
ExecStartPre=/bin/sh -c 'sudo systemctl restart bluetooth; sleep 5'
ExecStart=/home/pi/bms/bin/python /home/pi/bms_daemon.py
# Restart no matter how the process exits — a crash, a self-exit via
# sys.exit on stall detection, whatever. 15 s hold-off is plenty for
# the ExecStartPre bluetooth restart to settle.
Restart=always
RestartSec=15
# Kill after 30 s if graceful shutdown gets stuck (BluBleZ hangs on
# GATT teardown occasionally).
TimeoutStopSec=30
StandardOutput=append:/var/log/bms_daemon.log
StandardError=append:/var/log/bms_daemon.log

[Install]
WantedBy=multi-user.target
